# com.wooliesbot.automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wooliesbot.automation</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/danielkuschmierz/Woolies Script/.venv/bin/python</string>
        <string>/Users/danielkuschmierz/Woolies Script/chef_os.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/danielkuschmierz/Woolies Script</string>
    <key>StandardErrorPath</key>
    <string>/Users/danielkuschmierz/Woolies Script/logs/worker_error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/danielkuschmierz/Woolies Script/logs/worker_output.log</string>
    <key>RunAtLoad</key>
    <true/>
    <!-- Restart chef_os if it crashes; ThrottleInterval avoids tight restart loops -->
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
