# com.chefsos.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chefsos</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/danielkuschmierz/Woolies Script/.venv/bin/python</string>
        <string>chef_os.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/danielkuschmierz/Woolies Script</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/danielkuschmierz/Woolies Script/chef_os_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/danielkuschmierz/Woolies Script/chef_os_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>TELEGRAM_TOKEN</key>
        <string>YOUR_TOKEN_HERE</string>
        <key>TELEGRAM_CHAT_ID</key>
        <string>YOUR_CHAT_ID_HERE</string>
    </dict>
</dict>
</plist>
